# Assignment04/42-keyboard.rules
ACTION=="add", SUBSYSTEM=="usb", RUN+="/sbin/insmod /home/lucas/little-penguin/Assignment04/super-cool-keyboard-module.ko"
